# config/llama-embed.service
[Unit]
Description=llama.cpp Embedding Server (all-MiniLM-L6-v2)
After=network.target

[Service]
Type=simple
User=yantrik
Group=yantrik
ExecStart=/usr/local/bin/llama-server -m /opt/models/all-MiniLM-L6-v2-f16.gguf --host 127.0.0.1 --port 8082 --embedding --pooling mean -c 512
Restart=on-failure
RestartSec=5
MemoryMax=256M

[Install]
WantedBy=multi-user.target
